# llm_daily.cron
SHELL=/bin/zsh
0 23 * * * /usr/local/bin/node /Users/rahul/Desktop/React_app/llm_app/daily_llm_push.js >> /tmp/llm_daily.log 2>> /tmp/llm_daily.err
